# systemd/cmake-maintenance.timer
[Unit]
Description=CMake Maintenance Routine

[Timer]
OnCalendar=Daily
Persistent=true

[Install]
WantedBy=timers.target
