# systemd/llvm-maintenance.timer
[Unit]
Description=LLVM Maintenance Routine

[Timer]
OnCalendar=Daily
Persistent=true

[Install]
WantedBy=timers.target
